# crontab
# Run every day at 2:45 PM Eastern Standard Time
15 14 * * * python3 main.py
