# crontab
# Run every day at 5 AM Eastern Standard Time
0 10 * * * /usr/local/bin/python3.6 main.py
