# crontab
# Run every day at 5 AM Eastern Standard Time
0 10 * * * python3 main.py
