# crontab
# Run every day at 2:45 PM Eastern Standard Time
11 21 * * * python3 main.py
