# crontab
# Run every day at 4am 
0 4 * * * python main.py