# crontab
# Run every day at 5:15 AM Eastern Standard Time
15 10 * * * python3 main.py