# crontab
# Run every day at 4am 
0 4 * * * alembic upgrade head && python fbo.py