# crontab
# Run every day at 5 AM and 5PM Eastern Standard Time
5  10 * * * python3 main.py
5  22 * * * python3 main.py
