# crontab
# Run every day at 4am 
0 9 * * * python main.py