# crontab
# Run every day at 3:15 PM 
15 18 * * * python main.py