# crontab
# Run every day at 11:15 AM Eastern Standard Time
15 16 * * * python3 main.py
