# crontab
# Run every day at 5 AM and 5PM Eastern Standard Time
0  10 * * * python3 main.py
0  22 * * * python3 main.py
