# crontab
# Run every day at 4am 
0 6 * * * python main.py