# crontab
# Run every day at 4am 
0 4 * * * /usr/bin/env python3 /home/vcap/app/code/fbo.py